# roles/xfce-desktop/files/autostart/xscreensaver.desktop
[Desktop Entry]
Type=Application
Name=Screensaver
Exec=xscreensaver -no-splash
Hidden=true
X-GNOME-Autostart-enabled=false
